# deploy/bookstore.service
[Unit]
Description=Bookstore Console App
After=network.target

[Service]
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/bookstore
ExecStart=/usr/bin/java -Xms128m -Xmx512m -jar /opt/bookstore/bookstore.jar
Restart=on-failure
RestartSec=5
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=bookstore

[Install]
WantedBy=multi-user.target
